# plists/com.aim.alertd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/yangzs</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/yangzs/.local/bin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.aim.alertd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3.13</string>
		<string>/Users/yangzs/.aim/bin/alertd.py</string>
		<string>--persist-events</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/yangzs/.aim/logs/alertd.log</string>
	<key>StandardOutPath</key>
	<string>/Users/yangzs/.aim/logs/alertd.log</string>
	<key>SuccessfulExit</key>
	<false/>
	<key>WorkingDirectory</key>
	<string>/Users/yangzs/.aim/bin</string>
</dict>
</plist>
