# agents/ZS0002/com.aim.nats-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aim.nats-agent.ZS0002</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/yangzs/.aim/agents/ZS0002/nats-agent.py</string>
        <string>--agent-id</string>
        <string>ZS0002</string>
        <string>--server</string>
        <string>nats://127.0.0.1:4222</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>ThrottleInterval</key>
        <integer>5</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/yangzs/.aim/agents/ZS0002</string>
    <key>StandardOutPath</key>
    <string>/Users/yangzs/.aim/agents/ZS0002/logs/agent.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yangzs/.aim/agents/ZS0002/logs/agent.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>AIM_AGENT_ID</key>
        <string>ZS0002</string>
        <key>AIM_NATS_SERVER</key>
        <string>nats://127.0.0.1:4222</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
